# pi/motionsensor.service
[Unit]
Description=Face Recognition Software
After=multi-user.target

[Service]
Type=idle
WorkingDirectory=/home/pi/Documents/dataiker-recognition/pi
ExecStart=python3 /home/pi/Documents/dataiker-recognition/pi/motionsensor.py
User=pi

[Install]
WantedBy=multi-user.target
